# acfshell/service/xyz.openbmc_project.acfshell.service
[Unit]
Description=Script Runner Service

[Service]
ExecStart=/usr/bin/acfshell 
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target